# provisioning/cfg/supervisord.conf
[program:kcsrv]
command=/vagrant/start.sh
directory=/vagrant
user=vagrant

[program:kcproxy]
command=/vagrant/KCProxy/start.sh
directory=/vagrant/KCProxy
user=vagrant

