# provisioning/init/emperor.conf
description "uWSGI Emperor"
start on runlevel [2345]
stop on runlevel [06]

exec uwsgi --emperor /etc/uwsgi/apps-enabled --vassals-include /etc/uwsgi/vassals-default.ini --logto /var/log/uwsgi.log
